# docscheck.service
[Unit]
Description=DocsCheck Passport Verification Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/DocsCheck
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/DocsCheck/venv/bin/python3 bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
